# netstart.service
[Unit]
Description=fix network stuff

[Service]
ExecStart=/home/alex/shn/netstart.sh
Type=idle
User=root

[Install]
WantedBy=multi-user.target
