# archarm/sh/net-start.service
[Unit]
Description=fix network stuff

[Service]
ExecStart=/home/alex/sh/net-start.sh
Type=idle
User=root

[Install]
WantedBy=multi-user.target
