# bash/services/turnitoff.service
[Unit]
Description=Turn off ledstrip and amp.
Before=sleep.target

[Service]
User=
Type=
Environment=
ExecStart=
ExecStartPost=

[Install]
WantedBy=